# deploy/systemd/mon-agent.service
[Unit]
Description=mon agent (read-only host monitor)
Documentation=https://github.com/pr0ph37/mon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=monagent
Group=monagent
# Optional: grant collector access by adding monagent to these groups
# during package install. The agent degrades gracefully if absent.
#  - docker / podman : workload collector
#  - libvirt         : KVM/LXC inventory via virsh
#  - systemd-journal : sshd auth events via journalctl
#  - fail2ban        : fail2ban-client UNIX socket
SupplementaryGroups=docker libvirt systemd-journal fail2ban
ExecStart=/usr/local/bin/mon-agent --config=/etc/mon-agent/config.yaml

# Filesystem
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/mon-agent /var/log/mon-agent
# /var/log/wtmp, /var/log/btmp for last/lastb fallback when journal absent.
ReadOnlyPaths=/var/log/wtmp /var/log/btmp /etc/passwd /etc/group
PrivateTmp=yes
PrivateDevices=yes
ProtectHostname=yes
ProtectClock=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectProc=invisible
# ProcSubset=pid removed: hides /proc/stat, /proc/net/dev, /proc/meminfo
# which the agent's net/cpu/mem collectors need. ProtectProc=invisible
# already hides other users' process info, which is the real concern.

# Privileges
NoNewPrivileges=yes
# CAP_NET_ADMIN: nft list ruleset / iptables-save read live netfilter state.
# CAP_DAC_READ_SEARCH: read /var/log/wtmp, /var/log/btmp without group ownership.
CapabilityBoundingSet=CAP_NET_ADMIN CAP_DAC_READ_SEARCH
AmbientCapabilities=CAP_NET_ADMIN CAP_DAC_READ_SEARCH
LockPersonality=yes
MemoryDenyWriteExecute=yes
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX AF_NETLINK
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources @mount @swap @reboot @raw-io @debug @cpu-emulation @obsolete

UMask=0077
Restart=on-failure
RestartSec=5s
TimeoutStopSec=15s
# WatchdogSec removed: requires Type=notify and the agent calling
# sd_notify(WATCHDOG=1) on every successful tick; the current Type=simple
# agent doesn't implement that, so leaving WatchdogSec=60s caused
# systemd to kill it every 60s (Failed with result 'watchdog').

[Install]
WantedBy=multi-user.target
